# cache-cleaner.cron
*/5 * * * * curl -X PURGE https://coldnsteel.github.io/comedy-lounge-jokebox/
*/5 * * * * curl -X PURGE https://coldnsteel.github.io/MyWebsite/public/
*/5 * * * * curl -X PURGE https://coldnsteel.github.io/HackerWatch/public/
*/5 * * * * curl -X PURGE https://coldnsteel.github.io/KOZMIC-KASINO-/
